# roles/uwsgi/templates/uwsgi.conf
description "uWSGI Emperor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/{{ app.user }}/uwsgi/bin/uwsgi --emperor /etc/uwsgi/vassals/
